# ansible/templates/nextjs.service.j2
[Unit]
Description=Next.js E-Commerce Application
After=network.target

[Service]
Type=simple
User={{ app_user }}
WorkingDirectory={{ app_dir }}
Environment="NODE_ENV=production"
Environment="PORT={{ app_port }}"
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ app_name }}

[Install]
WantedBy=multi-user.target
